# puppet/modules/fx/files/etc/cron.d/forexite_download_data
MAILTO=joaocosta@zonalivre.org
TRADER_HOME=/home/joao/src/Finance--HostedTrader

8 7 * * * joao $TRADER_HOME/data/forexite/update 60
